# spoofdpi.service
[Unit]
Description=dpi spoofer
After=network.target

[Service]
ExecStart=%h/spoofdpi
Type=simple
Restart=always

[Install]
WantedBy=default.target
